# tg-news-bot.service
[Unit]
Description=Telegram News Bot (Инфобиз/Маркетинг)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/tg_news_bot
ExecStart=/home/ubuntu/tg_news_bot/venv/bin/python main.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Автоматический перезапуск при падении
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
